# templates/systemd/mongod_thp.service.erb
[Unit]
Description=Mongos_<%= @mongod_instance %>_THP_Settings
Wants=network.target
After=network.target
Before=mongod_<%= @mongod_instance %>.service

[Service]
Type=oneshot
RemainAfterExit=true
<% if @mongod_deactivate_transparent_hugepage -%>
<%# implement mdiag.sh best-practice %>
ExecStartPre=/bin/bash -c 'test -f /sys/kernel/mm/transparent_hugepage/khugepaged/defrag && echo 0 > /sys/kernel/mm/transparent_hugepage/khugepaged/defrag'
ExecStartPre=/bin/bash -c 'test -f /sys/kernel/mm/transparent_hugepage/defrag && echo never > /sys/kernel/mm/transparent_hugepage/defrag'
ExecStartPre=/bin/bash -c 'test -f /sys/kernel/mm/transparent_hugepage/enabled && echo never > /sys/kernel/mm/transparent_hugepage/enabled'
<% end -%>
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
